# server/misc/ros2-node-manager.target
[Unit]
Description=Group of ROS2 nodes that are managed by the ROS2 node manager

